# jibri/rootfs/etc/s6-overlay/s6-rc.d/jibri/run
#!/command/execlineb
with-contenv
s6-setuidgid jibri
/etc/s6-overlay/scripts/run-jibri
