# prosody/rootfs/etc/services.d/prosody/run
#!/command/with-contenv bash
exec prosody --config /config/prosody.cfg.lua -F
