# prosody/rootfs/etc/services.d/prosody/run
#!/usr/bin/with-contenv bash
exec prosody --config /config/prosody.cfg.lua

