# jibri/rootfs/etc/s6-overlay/s6-rc.d/40-jibri/finish
#!/command/with-contenv bash

exec /etc/s6-overlay/scripts/jibri-finish "$@"
