# web/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb
/etc/s6-overlay/scripts/run-nginx
