# jvb/rootfs/etc/s6-overlay/s6-rc.d/jvb/run
#!/command/execlineb
with-contenv
s6-setuidgid jvb
/etc/s6-overlay/scripts/run-jvb
