# sip-jibri/rootfs/etc/services.d/20-icewm/run
#!/usr/bin/with-contenv bash

DAEMON="/usr/bin/icewm-session"
exec s6-setuidgid jibri /bin/bash -c "exec $DAEMON"

