# prosody/rootfs/etc/services.d/10-saslauthd/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid root saslauthd -a ldap -O /etc/saslauthd.conf -c -m /var/run/saslauthd -n 5 -d
