# prosody/rootfs/etc/services.d/10-saslauthd/run
#!/command/with-contenv bash
exec saslauthd -a ldap -O /config/saslauthd/saslauthd.conf -c -m /var/run/saslauthd -n 5 -d
