# web/rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

if [[ $DISABLE_HTTPS -ne 1 ]]; then
    if [[ $ENABLE_LETSENCRYPT -eq 1 ]]; then
      exec cron -f
    fi
fi
