# web/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
nginx -c /config/nginx/nginx.conf

