# prosody/rootfs/etc/s6-overlay/s6-rc.d/saslauthd/run
#!/command/execlineb
s6-setuidgid root
/etc/s6-overlay/scripts/run-saslauthd
