# prosody/rootfs/etc/s6-overlay/s6-rc.d/prosody/run
#!/command/execlineb
s6-setuidgid prosody
/etc/s6-overlay/scripts/run-prosody
