# jibri/rootfs/etc/services.d/30-pulse/run
#!/command/with-contenv bash

HOME=/home/jibri
exec /bin/bash -c "exec /usr/bin/pulseaudio"
