# jibri/rootfs/etc/services.d/30-pulse/run
#!/command/with-contenv bash

HOME=/jibri
exec s6-setuidgid /bin/bash -c "exec /usr/bin/pulseaudio"
