# sip-jibri/rootfs/etc/services.d/40-icewm-display-1/run
#!/usr/bin/with-contenv bash

DAEMON="/usr/bin/icewm-session -d ${DISPLAY1}"
exec s6-setuidgid jibri /bin/bash -c "exec $DAEMON"

