# jicofo/rootfs/etc/s6-overlay/s6-rc.d/jicofo/run
#!/command/execlineb
with-contenv
s6-setuidgid jicofo
/etc/s6-overlay/scripts/run-jicofo
