# web/rootfs/etc/s6-overlay/s6-rc.d/cron/run
#!/command/execlineb
with-contenv
/etc/s6-overlay/scripts/run-cron
