# sip-jibri/rootfs/etc/services.d/50-pulse/run
#!/usr/bin/with-contenv bash

HOME=/home/jibri
exec s6-setuidgid jibri /bin/bash -c "exec /usr/bin/pulseaudio"
